# packages/platform/raspberrypi/otto-boot/files/otto-bq27510.service
[Unit]
Before=otto-boot.service

[Service]
ExecStart=/bin/bash -c "echo \"bq27500 0x55\" > /sys/class/i2c-dev/i2c-1/device/new_device"

[Install]
WantedBy=multi-user.target