# packages/platform/raspberrypi/otto-wifi/files/dnsmasq.service
[Unit]
Description=A lightweight DHCP and caching DNS server
After=network.target otto-wifi.service
Documentation=man:dnsmasq(8)

[Service]
Type=dbus
BusName=uk.org.thekelleys.dnsmasq
ExecStartPre=/usr/sbin/dnsmasq --test
ExecStartPre=/sbin/ip addr add 192.168.12.1/24 broadcast 192.168.12.255 dev wlan1
ExecStart=/usr/sbin/dnsmasq -k --enable-dbus -x /var/run/dnsmasq.pid
ExecReload=/bin/kill -HUP $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
