# packages/platform/raspberrypi/otto-hostname/files/otto-hostname.service
[Unit]
Description=OTTO bootup hostname script
After=syslog.target auditd.service

[Service]
ExecStart=/bin/bash -c "/bin/hostname otto-$(cat /sys/class/net/eth0/address | head -1 | sed 's/\://g')"
Restart=no

[Install]
WantedBy=multi-user.target