# packages/platform/raspberrypi/otto-shutdown/files/otto-shutdown.service
[Unit]
Description=Pull pin 4 low on shutdown
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/otto-poweroff
Type=oneshot

